# BambiCTF_2021/ringring/RotatorCronJob/log_rotator.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * /bin/sh /rotate_files.sh >> /var/log/script.log
